# root/etc/s6-overlay/s6-rc.d/init-kavita-config/up
/etc/s6-overlay/s6-rc.d/init-kavita-config/run
